# module/config/crontab
*/3 * * * *     bash /home/${username}/zomboid/backup_server.sh > /home/${username}/zomboid/backup_server.log
